# openwrt/openwrt-config/crontab
* * * * * /home/scripts/update_GNSM
* * * * * /home/scripts/update_ping

